# proxmox/templates/pve-metrics.service.j2
# pve-metrics.service.j2
[Unit]
Description=Proxmox Metrics Exporter
Documentation=https://prometheus.io/docs/instrumenting/exporters/
After=network.target pve-cluster.service
Wants=network.target

[Service]
Type=simple
User={{ otel_user }}
Group={{ otel_group }}
# CAP_SYS_ADMIN is added to allow access for commands like nvme and sensors to read hardware data
AmbientCapabilities=CAP_SYS_ADMIN
ExecStart=/bin/bash {{ script_path }}
Restart=always
RestartSec=30
TimeoutStopSec=20

# Security settings
NoNewPrivileges=yes
ProtectSystem=strict
ProtectHome=yes
ReadWritePaths=/var/lib/otelcol-contrib
# Added /dev to ReadOnlyPaths to allow smartctl to access device information (though our script uses nvme-cli directly)
ReadOnlyPaths=/sys /proc /etc/pve /dev
PrivateTmp=yes
ProtectHostname=yes
ProtectClock=yes
ProtectKernelTunables=yes
ProtectKernelModules=yes
ProtectKernelLogs=yes
RestrictRealtime=yes
RestrictSUIDSGID=yes
RemoveIPC=yes

# Resource limits
MemoryMax=512M
CPUQuota=25%

# Environment (adjust if your script uses Python environments or other specific settings)
Environment=PYTHONUNBUFFERED=1
Environment=FLASK_DEBUG=0 # Updated for Flask 2.3+
# Environment=FLASK_ENV=production # Old, deprecated setting - commented out as per your original

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=pve-metrics-exporter

[Install]
WantedBy=multi-user.target
