# proxmox/templates/pve-otelcol-contrib.service.j2
# /etc/systemd/system/pve-metrics-exporter.service
# pve-metrics.service.j2
[Unit]
Description=Proxmox Metrics Exporter
Documentation=https://prometheus.io/docs/instrumenting/exporters/
After=network.target pve-cluster.service
Wants=network.target

[Service]
Type=simple
User=otelcol
Group=otelcol
ExecStart=/usr/bin/python3 /opt/otelcol-contrib/pve-metrics-exporter.py
Restart=always
RestartSec=5
TimeoutStopSec=20

# Security settings
NoNewPrivileges=yes
ProtectSystem=strict
ProtectHome=yes
ReadWritePaths=/var/lib/otelcol-contrib
# Crucial change for smartctl: allow read access to /dev
ReadOnlyPaths=/sys /proc /etc/pve /dev
PrivateTmp=yes
ProtectHostname=yes
ProtectClock=yes
ProtectKernelTunables=yes
ProtectKernelModules=yes
ProtectKernelLogs=yes
RestrictRealtime=yes
RestrictSUIDSGID=yes
RemoveIPC=yes

# Resource limits
MemoryMax=512M
CPUQuota=25%

# Environment
Environment=PYTHONUNBUFFERED=1
Environment=FLASK_ENV=production

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=pve-metrics-exporter

[Install]
WantedBy=multi-user.target
