# pve/roles/metrics-exporter/templates/metrics-exporter.service.j2
[Unit]
Description=Proxmox Metrics Exporter
Documentation=https://prometheus.io/docs/instrumenting/exporters/
After=network.target pve-cluster.service
Wants=network.target

[Service]
Type=simple
User={{ metrics_user }}
Group={{ metrics_group }}
# CAP_SYS_ADMIN is added to allow access for commands like nvme and sensors to read hardware data
AmbientCapabilities=CAP_SYS_ADMIN
ExecStart=/usr/bin/python3 {{ metrics_exporter_app_path }}
Restart=always
RestartSec=5
TimeoutStopSec=20

# Security settings
NoNewPrivileges=no  # Changed from 'yes' to 'no' to allow capability usage
ProtectSystem=strict
ProtectHome=yes
ReadWritePaths=/var/lib/otelcol-contrib
# Added /dev to ReadOnlyPaths to allow smartctl/nvme to access device information
ReadOnlyPaths=/sys /proc /etc/pve /dev
PrivateTmp=yes
ProtectHostname=yes
ProtectClock=yes
ProtectKernelTunables=yes
ProtectKernelModules=yes
ProtectKernelLogs=yes
RestrictRealtime=yes
RestrictSUIDSGID=yes
RemoveIPC=yes

# Resource limits
MemoryMax=512M
CPUQuota=25%

# Environment
Environment=PYTHONUNBUFFERED=1
Environment=FLASK_DEBUG=0
Environment=NO_COLOR=1
Environment=TERM=dumb

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=pve-metrics-exporter

[Install]
WantedBy=multi-user.target
